# var-opt-www/demoapp/demoapp_supervisor.ini
[program:demoapp]
command=/usr/local/bin/uwsgi
  --ini /var/opt/www/demoapp/demoapp_uwsgi.ini
directory=/var/opt/www/demoapp
autostart=true
autorestart=true